# roles/minio/templates/minio.service.j2
[Unit]
Description=MinIO Object Storage
After=network.target

[Service]
User=root
Group=root
EnvironmentFile=/etc/minio.env
ExecStart=/usr/local/bin/minio server {{ minio_data_dir }} \
  --console-address :{{ minio_console_port }}
Restart=always
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target